# profiles/oneplus-fajita/overlay/usr/lib/systemd/system/iio-sensor-proxy.service.d/oneplus6t-delay.conf
[Service]
# PMOS starts iio-sensor-proxy after hexagonrpcd has had time to settle.
ExecStartPre=/usr/bin/sleep 10
